# recipes-liion-ab-bms/liion-hmi-app-init/files/liion_hmi_autorun.service
[Unit]
After=networking.service

[Service]
ExecStart=/home/liion/liion_hmi_autorun.sh

[Install]
WantedBy=default.target